# rootfs/etc/services.d/node-bob/run
#!/bin/bash

echo "* Starting node-bob ..."
/usr/bin/bitcoind -daemon -datadir=/data/node-bob

s6-svc -d /var/run/s6/services/node-bob
