# rootfs/etc/services.d/node-alice/run
#!/bin/bash

echo "* Starting node-alice ..."
/usr/bin/bitcoind -daemon -datadir=/data/node-alice

s6-svc -d /var/run/s6/services/node-alice

